# recipes-debugging/notify-spammer/files/notify-spammer.service
[Unit]
Description=Watchdog spammer

[Service]
Type=notify
ExecStart=/usr/sbin/notify-spammer
WatchdogSec=10

[Install]
WantedBy=multi-user.target